# data/passive-ddns.service.default
[Unit]
Description=Passive DDNS service
Wants=passive-ddns.timer
After=network.target

[Service]
Type=oneshot
DynamicUser=true
Environment=RUST_LOG=info
WorkingDirectory=/opt/passive-DDNS
ExecStart=/opt/passive-DDNS/target/release/passive-ddns --systemd

[Install]
WantedBy=multi-user.target

