# data/passive-ddns.service.default
[Unit]
Description=Passive DDNS service
After=network.target

[Service]
Type=simple
DynamicUser=true
Environment="RUST_LOG=info"
Restart=on-failure
RestartSec=10s
WorkingDirectory=/opt/passive-DDNS
ExecStart=/opt/passive-DDNS/target/release/passive-ddns --systemd

[Install]
WantedBy=multi-user.target
