# data/passive-ddns.service.default
[Unit]
Description=Passive DDNS service
After=network.target

[Service]
Type=simple
User=nobody
Restart=on-failure
RestartSec=5s
WorkingDirectory=/opt/passive-DDNS
ExecStart=/usr/bin/python3 /opt/passive-DDNS/hostkerddns.py
ExecReload=/bin/kill -10 $MAINPID
ExecStop=/usr/bin/python3 /opt/passive-DDNS/absddns.py stop $MAINPID
TimeoutStopSec=5s

[Install]
WantedBy=multi-user.target

