# data/passive-ddns.timer.default
[Unit]
Description=Auto start Passive DDNS service
Wants=passive-ddns.service

[Timer]
Unit=passive-ddns.service
OnCalendar=*:0/10

[Install]
WantedBy=timers.target